# files/apt-daily-upgrade.override.timer
[Timer]
OnCalendar=
OnCalendar=*-*-* 04:00
RandomizedDelaySec=45m
Persistent=true